# etc/xdg/autostart/gooroom-keycode.desktop
[Desktop Entry]
Type=Application
Name=gooroom keycode setting
Exec=gooroom-keycode
NoDisplay=true
